# scripts/kubmix-hctg-worker.service
[Unit]
Description=Betano favorito screenshot worker (Kubmix sa-east-1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/betano-monitor/scripts
Environment=NODE_ENV=production
Environment=HCTG_WORKER_SOURCE=kubmix-worker
Environment=HCTG_HTML_SOURCE=html-dom-kubmix
Environment=HCTG_HEADLESS=0
Environment=HCTG_SCRAPE=0
Environment=FAVORITO_SCREENSHOT=1
Environment=FAVORITO_POLL_SEC=20
Environment=HCTG_CHROME_PROFILE=/home/admin/betano-monitor/.chrome-hctg-profile
Environment=HCTG_POLL_SEC=90
Environment=HCTG_MAX_PER_CYCLE=4
Environment=HCTG_BROWSER_RESTART_MIN=30
EnvironmentFile=-/home/admin/betano-monitor/.env
# So prints do odd inicial (favorito). HCTG Total de Gols desligado (HCTG_SCRAPE=0).
ExecStart=/usr/bin/xvfb-run -a -s "-screen 0 1920x1080x24" /usr/bin/node hctg-worker.mjs
Restart=always
RestartSec=20
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
